# madhatter-sync.service
[Unit]
Description=Madhatter Drop Sync Daemon
After=network.target

[Service]
ExecStart=/usr/bin/madhatter-sync
Restart=always
RestartSec=10
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
